# etc/init.d/set_keyboard_rate
#!/sbin/openrc-run
name="set-keyboard-rate"

description="sets keyboard rate"




start() {
    ebegin "set keyboard rate"
    kbdrate -r 100 -d 130
    eend $?
}
